# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "starting cron"

#fix link-count, as cron is being a pain, and docker is making hardlink count >0 (very high)
touch /etc/crontab /etc/cron.*/*

exec cron -L 8 -f
