# rootfs/etc/services.d/klaxon/finish
#!/usr/bin/execlineb -S0

echo "klaxon exiting"
s6-svscanctl -t /var/run/s6/services