# rootfs/etc/services.d/klaxon/run
#!/usr/bin/with-contenv bash

echo "starting klaxon"
cd /usr/src/app
bundle exec puma -C "config/puma.rb"
